# services/autodeploy.timer
[Unit]
Description=Schedule autodeployment  execution

[Timer]
OnUnitActiveSec=300s
OnBootSec=300s

[Install]
WantedBy=timers.target
